# listenr-web@.service
[Unit]
Description=Listenr Web Server
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/Code/listenr
Environment=PATH=/home/%i/Code/listenr/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/%i/Code/listenr/venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 2 --timeout 120 webapp:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target